# services/prepperpi-admin/prepperpi-admin.service
[Unit]
Description=PrepperPi — admin console (FastAPI/uvicorn)
Documentation=https://github.com/jmarler/prepperpi
After=network.target

[Service]
Type=simple
User=prepperpi-admin
Group=prepperpi-admin
WorkingDirectory=/opt/prepperpi/services/prepperpi-admin/app
ExecStart=/usr/bin/python3 -m uvicorn main:app \
  --host 127.0.0.1 \
  --port 8090 \
  --no-access-log \
  --log-level info
Restart=on-failure
RestartSec=5s

# Sandboxing. The daemon never writes to disk except indirectly via
# the apply-network-config sudo wrapper; everything else is read-only.
NoNewPrivileges=false
# ^ false because we DO need to invoke sudo. The privilege bump is
# bounded by the sudoers rule (one specific script, NOPASSWD only).
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
# AF_NETLINK is required so the uplink detector can shell out to
# `ip -j route show default`. iproute2 uses netlink to query the kernel;
# without it, the call fails with "Cannot open netlink socket" and the
# admin home page falsely renders as Offline. AF_NETLINK only allows
# read-style queries here — the unit still has no CAP_NET_ADMIN, so it
# can't *write* anything to the kernel via netlink.
ReadOnlyPaths=/opt/prepperpi/services/prepperpi-admin
# /boot/firmware has to be writable IN THE NAMESPACE so the sudo'd
# wrapper can rewrite prepperpi.conf -- a read-only bind mount blocks
# even root. UNIX perms still keep the unprivileged uvicorn process
# (which has no write permission to the file) from touching it
# directly; only the wrapper, escalated via sudo to root, can write.
ReadWritePaths=/boot/firmware /srv/prepperpi/cache /srv/prepperpi/maps /srv/prepperpi/zim /var/lib/prepperpi/bundles /etc/prepperpi/bundles /var/lib/prepperpi/updates /var/lib/prepperpi/backup
# /srv/prepperpi/cache for the kiwix-catalog.json cache.
# /srv/prepperpi/maps for the region-uninstall path: the admin
# console deletes <region>.mbtiles directly (no privileged worker —
# the maps dir is owned by prepperpi-admin so plain unlink suffices).
# /var/lib/prepperpi/bundles caches fetched remote bundle manifests
# and per-install state. /etc/prepperpi/bundles holds the editable
# sources.json (admin can add community URLs).
# /var/lib/prepperpi/updates holds state.json (latest update-availability
# snapshot) and pins.json (per-item update pins).
#
# /srv/prepperpi/zim is added so the admin can unlink ZIMs (e.g. old
# versions left behind after a side-by-side update). aria2c (running
# as `prepperpi`, which OWNS the directory) still does all the
# *creating*; admin only deletes. UNIX perms gate further: setup.sh
# adds prepperpi-admin to the `prepperpi` group and chmods the dir to
# 775 so group members can unlink files in it.
#
# We DO NOT add /srv/prepperpi/user-usb here. Adding it creates a bind
# mount in admin's mount namespace, which silently breaks the USB
# write-toggle wrapper: sudo'd `mount -o remount,rw` would target the
# bind mount, not the underlying USB device, so the host's /proc/mounts
# stays ro and writes still fail.
LockPersonality=true
RestrictNamespaces=mnt
# Was `true` (deny everything), but `RestrictNamespaces=true` installs
# a seccomp filter that blocks setns(2), which prevents the sudo'd
# apply-storage-action wrapper from doing
#     nsenter --target 1 --mount -- mount -o remount,rw,...
# to escape into the host's mount namespace. Without that escape, the
# remount only mutates admin's private namespace view; the host's
# /proc/mounts stays read-only and aria2c's writes fail with EACCES.
# Allowing only `mnt` keeps the other namespace types denied.
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
